# configs/supervisor/oniichan.conf
[program:oniichan]
directory = /var/www/oniichan/
command = /var/www/oniichan/prod/bin/python /var/www/oniichan/prod/bin/gunicorn oniichan:app
username = oniichan